# linux/systemd/socket/demo2/echo2@.service
[Service]
ExecStart=/usr/local/bin/echo2
StandardInput=socket
StandardOutput=socket
